# crontab
PATH=/usr/local/bin:/usr/bin:/bin

# Run trading bot at 9 AM daily (Europe/Bucharest time)
0 9 * * * cd /app && /usr/local/bin/python3 automated_runner.py >> /var/log/cron.log 2>&1

# Update performance metrics at 10 AM daily (after markets settle)
0 10 * * * cd /app && /usr/local/bin/python3 track_performance.py >> /var/log/cron.log 2>&1

# Empty line required at end of crontab file
